# deploy/systemd/nl-table-healthcheck.timer
[Unit]
Description=Run NL Table production health check every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
